# ingress/dist/vegmod.crontab
# Ingress every 3 minutes
*/3 * * * * /usr/local/bin/vegmod ingress >> /app/logs/praw.log 2>&1
